# Build/supervisor.conf
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=900
redirect_stderr=true
redirect_stdout=true
autostart=true
autorestart=true
startretries=10
startsecs=4
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
user=root
username=www-data

[program:gunicorn]
command=/app/ntss_app/bin/gunicorn wsgi:application -c /etc/gunicorn.conf.py --reload
directory=/app/ntss_www
redirect_stderr=true
redirect_stdout=true
autostart=true
autorestart=true
startretries=10
startsecs=4
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
user=nobody

[program:postfix]
command=/usr/sbin/postfix -c /etc/postfix start-fg
autostart=true
autorestart=true
startretries=10
startsecs=4
redirect_stderr=true
redirect_stdout=true
startsecs=0