# Build/supervisor.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/apache2ctl -DFOREGROUND
redirect_stderr=true
redirect_stdout=true
autostart=true
autorestart=true
startretries=1
startsecs=1
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
user=root
